# systemd/pistreamer.service
[Unit]
Description=PiStreamer Web and YouTube streaming service
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
SupplementaryGroups=video audio
WorkingDirectory=/opt/pistreamer
Environment=PISTREAMER_CONFIG=/etc/pistreamer/config.yaml
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/pistreamer/.venv/bin/python /opt/pistreamer/run.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
