# systemd/pistreamer.service
[Unit]
Description=PiStreamer Web and streaming service
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
SupplementaryGroups=video audio
WorkingDirectory=/opt/pistreamer
Environment=PISTREAMER_CONFIG=/etc/pistreamer/config.yaml
Environment=PYTHONUNBUFFERED=1
Environment=TZ=Europe/Berlin
RuntimeDirectory=pistreamer
RuntimeDirectoryMode=0775
ExecStart=/opt/pistreamer/.venv/bin/python /opt/pistreamer/run.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
