# systemd/pistreamer-update.path
[Unit]
Description=Watch for PiStreamer web update requests

[Path]
PathExists=/run/pistreamer/update.request
Unit=pistreamer-update.service

[Install]
WantedBy=multi-user.target
